# every2minutes.cron
*/2 * * * * /Users/jamiepatel/Desktop/projects/whos_in/local_scanner.sh